# docker_files/etc/service/contactqr/run
#!/bin/bash

set -e

exec /opt/release/contactqr-${VERSION}-linux-amd64 -path /opt/release/ui