# docker_files/etc/service/contactqr/run
#!/bin/bash

set -e

exec /opt/contactqr/contactqr-${VERSION}-linux-amd64 -port ${PORT} -path /opt/contactqr/ui